# system/usr/lib/systemd/system/bootc-upgrade.timer
[Unit]
Description=Bootc automatic upgrade timer

[Timer]
OnCalendar=*-*-* 14:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
